# deploy/systemd/aula-hacky-poll.service
[Unit]
Description=Poll and sync AULA keyboard RTC

[Service]
Type=oneshot
ExecStart=/usr/local/bin/aula-hacky-poll.sh
